# supervisor.conf
[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[program:worker]
command=celery -A lolsite worker --loglevel=INFO

[program:web]
command=gunicorn --bind :8000 --workers 2 lolsite.wsgi:application

[program:redis]
command=redis-server
